# deploy/m20ctrl.service
[Unit]
Description=M20 mobile web controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user
Group=user
WorkingDirectory=/home/user/m20ctrl
Environment=MAP_ROOT=/home/user/m20-fastlio/maps
Environment=HTTP_PORT=8000
Environment=WS_PORT=8080
Environment=M20_HOST=10.21.31.103
Environment=M20_PORT=30001
Environment=NAV_HOST=127.0.0.1
Environment=NAV_PORT=30011
Environment=NAV_ROUTE_MODE=custom
Environment=NAV_VIZ_HOST=127.0.0.1
Environment=NAV_VIZ_PORT=30013
Environment=NAV_VIZ_REGISTER_PORT=30012
ExecStart=/home/user/m20ctrl/start_nos_python.sh
Restart=on-failure
RestartSec=2
KillSignal=SIGTERM
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
